# osx/Library/LaunchAgents/com.teamdrive.teamdrive3.plist

                  <plist version="1.0">
                    <dict>
                      <key>Label</key>
					  <string>com.teamdrive.teamdrive3</string>
                      <key>Program</key>
                      <string>/Applications/TeamDrive/TeamDrive.app/Contents/MacOS/TeamDrive</string>
                      <key>RunAtLoad</key>
                      <true/>
					  <key>ProgramArguments</key>
					  <array>
						<string>autostart</string>
					  </array>
                    </dict>
                  </plist>
                  
                        